# OpenBSD-stuff/etc/rc.d/emacsd
#!/bin/ksh
#
# $OpenBSD: wesnothd.rc,v 1.2 2018/01/11 19:27:02 rpe Exp $

daemon="/usr/local/bin/emacs"
daemon_flags="--daemon"
daemon_user="greenfork"

. /etc/rc.d/rc.subr

rc_cmd $1
